# etc/systemd/system/geminibot.service
[Unit]
Description=Gemini Telegram Bot
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/geminibot
Environment=PATH=/opt/geminibot/venv/bin
ExecStart=/opt/geminibot/venv/bin/python main.py
Restart=always
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=15
KillSignal=SIGTERM
KillMode=mixed

[Install]
WantedBy=multi-user.target